# packaging/systemd/grlx-farmer.service
[Unit]
Description=grlx Farmer - Control Plane Service
Documentation=https://grlx.dev
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=farmer
Group=farmer
ExecStart=/usr/bin/grlx-farmer
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/cache/grlx/farmer /etc/grlx/pki/farmer
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Working directory and environment
WorkingDirectory=/var/cache/grlx/farmer
Environment=GRLX_CONFIG=/etc/grlx/farmer

[Install]
WantedBy=multi-user.target 